# private_dot_config/systemd/user/adhd-wall.service
[Unit]
Description=Set the wallpaper requested by the quickshell wallpaper widget

[Service]
Type=oneshot
ExecStart=%h/.local/bin/adhd-wall-exec.sh
