# private_dot_config/systemd/user/adhd-wall.path
[Unit]
Description=Watch for wallpaper-set requests from the quickshell wallpaper widget

[Path]
PathModified=%h/.cache/adhd/wall-request
Unit=adhd-wall.service

[Install]
WantedBy=default.target
