# private_dot_config/systemd/user/adhd-block.service
[Unit]
Description=Execute an ADHD focus-block start requested by the quickshell drawer

[Service]
Type=oneshot
ExecStart=%h/.local/bin/adhd-block-exec.sh
