# private_dot_config/systemd/user/adhd-prayer-times.service
[Unit]
Description=Regenerate ADHD-OS prayer times (offline, Hyderabad/Hanafi)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/adhd-prayer-times.sh
