# private_dot_config/systemd/user/adhd-task-cmd.service
[Unit]
Description=Execute taskwarrior commands queued by the quickshell task panel

[Service]
Type=oneshot
ExecStart=%h/.local/bin/adhd-task-cmd.sh
